# prosushimsg.service
[Unit]
Description=ProSushi Messenger - Corporate Chat Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/prosushimsg
# 🔧 Автоматически найти путь к dotnet
ExecStart=/usr/bin/env dotnet /var/www/prosushimsg/prosushimsg.dll
Restart=always
RestartSec=10

# 🔥 Переменные окружения
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# ⚠️ Лимиты памяти (для VPS с 512 МБ RAM)
MemoryMax=450M

# 📊 Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=prosushimsg

# 🔒 Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
